# infrastructure/docker/services/php/frontend/etc/service/php-fpm/run
#!/bin/sh

exec /usr/sbin/php-fpm${PHP_VERSION} -y /etc/php/${PHP_VERSION}/fpm/php-fpm.conf -O
