# infrastructure/docker/services/frontend/etc/service/php-fpm/run
#!/bin/sh

exec /usr/sbin/php-fpm7.4 -y /etc/php/7.4/fpm/php-fpm.conf -O
