# infrastructure/docker/services/frontend/etc/service/php-fpm/run
#!/bin/sh

exec php-fpm7 -y /etc/php7/php-fpm.conf -O
