# infrastructure/docker/services/cron/crontab
* * * * * php /home/app/application/bin/console app:cron >> /home/app/application/var/log/cron.log
